# root/etc/services.d/nebula/run
#!/usr/bin/with-contenv bash
umask 022

exec \
	s6-setuidgid abc nebula -config /config/config.yml
